# debian/postinst
#!/bin/bash
set -e
case "$1" in
  configure)
    # Install ASR runtimes via pip (not available in apt)
    echo "ByteCLI: Installing speech recognition engines..."
    echo "ByteCLI: This may take several minutes on first install."
    if [ -n "$SUDO_USER" ]; then
        su - "$SUDO_USER" -c \
          "/usr/bin/python3 -m pip install --user --extra-index-url https://download.pytorch.org/whl/cu121 faster-whisper openai-whisper 'funasr>=1.3.3' 'torchaudio==2.5.1' 'qwen-asr>=0.0.6' 'Pillow>=10'" \
          2>/dev/null || true
    else
        /usr/bin/python3 -m pip install --user --extra-index-url https://download.pytorch.org/whl/cu121 \
          faster-whisper openai-whisper 'funasr>=1.3.3' 'torchaudio==2.5.1' 'qwen-asr>=0.0.6' 'Pillow>=10' \
          2>/dev/null || true
    fi

    # Enable the user service globally (for all users, on next login)
    systemctl --global enable bytecli.service 2>/dev/null || true

    # Immediate start for the user who ran sudo
    if [ -n "$SUDO_USER" ]; then
        SUDO_UID=$(id -u "$SUDO_USER")
        if [ -d "/run/user/$SUDO_UID" ]; then
            su - "$SUDO_USER" -c \
              "systemctl --user daemon-reload && systemctl --user start bytecli.service" \
              2>/dev/null || true
        fi
    fi
    ;;
esac

exit 0
